# config/services.d/sync/run
#!/usr/bin/with-contenv sh

cd /data/
if [ -z ${SYNC_ONLY:+x} ]; then
	s6-pause;
else
	python sync.py
fi
