# config/services.d/frontend/run
#!/usr/bin/with-contenv sh

if [ -z ${SYNC_ONLY:+x} ];
then
	nginx -g "daemon off;"
else
	s6-pause;
fi
